# Api/supervisor/supervisord.conf
[supervisord]
nodaemon = true

[program:php-fpm]
command = /usr/sbin/php-fpm8.0 -F # Atualize para PHP 8.0

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"